# templates/jmxtrans.service.epp
[Unit]
Description=JMX Transformer - more than meets the eye
After=syslog.target network.target

[Service]
Type=forking
User=<%= $user %>
<% if $systemd_environment_file  { -%>
EnvironmentFile=<%= $systemd_environment_file %>
<% } -%>
<% if $working_directory { -%>
WorkingDirectory=<%= $working_directory %>
<% } -%>
PermissionsStartOnly=true
RuntimeDirectory=jmxtrans  
PIDFile=/var/run/jmxtrans/jmxtrans.pid
ExecStart=<%= $binary_path %> start

[Install]
WantedBy=multi-user.target
